# logging/cron-log-plotter/crontab
#0 * * * * ntpdate 0.arch.pool.ntp.org
45 12,18,23 * * * bash /home/drill/surface-unit/make-daily-drilllog-plot.sh

