# logging/cron-log-plotter/crontab
#0 * * * * ntpdate 0.arch.pool.ntp.org
59 11,15,19,23 * * * bash /home/drill/make-daily-drilllog-plot.sh



